# examples/xandikos.socket
[Unit]
Description=Xandikos socket

[Socket]
ListenStream=/run/xandikos.sock

[Install]
WantedBy=sockets.target
